# configs/tp6.service
[Unit]
Description = Tp6

[Service]
ExecStart = ../bin/api_contador
ExecStart = ../bin/api_gestor_usuario

[Install]
WantedBy = multi-user.target
